# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

; Sin estas dos secciones, `supervisorctl status` responde «.ini file does not include
; supervisorctl section» y no hay forma de ver, parar ni reiniciar un proceso desde
; dentro del contenedor. El día que algo se cae, es lo primero que uno quiere teclear.
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; El reloj del sistema. Sin esto, NADA de lo programado ocurre: los preventivos no
; se generan, el mes no se cierra, las alertas no se levantan y —lo más grave— el
; respaldo diario de la base de datos nunca corre. No falla: simplemente no pasa, y
; no hay una sola línea de log que lo diga. Estuvo ausente desde el primer día.
;
; `schedule:work` es un proceso vivo que dispara las tareas cada minuto; es el
; equivalente del cron del host, pero dentro del contenedor y sin cron.
[program:scheduler]
command=php artisan schedule:work
directory=/var/www/html
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=25

; Consume TODAS las colas que la aplicación usa. La lista tiene que coincidir con las
; que los jobs declaran en `onQueue()`: una cola que no esté aquí no falla, se queda
; con los trabajos dentro para siempre. Pasó con `maintenance` (el generador de
; preventivos) y `analytics` (el cierre mensual de KPIs) — meses sin ejecutarse, sin
; un error, mientras el tablero decía que los planes estaban activos. Hay un test que
; ahora compara esta lista contra los jobs y falla si vuelven a divergir.
;
; No se nombra la conexión a propósito: se resuelve QUEUE_CONNECTION del entorno, así
; que funciona igual con `database` o con `redis`. --max-time recicla el proceso cada
; hora para soltar memoria; autorestart lo levanta enseguida.
[program:queue-worker]
command=php artisan queue:work --queue=maintenance,automations,webhooks,analytics,audit,exports,default --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
autostart=true
autorestart=true
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
