# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
user=root

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

; Consumes every named queue used across the app (exports, webhooks,
; automations) plus the default queue. No connection is named on the command
; line on purpose — it resolves QUEUE_CONNECTION from the environment at
; runtime, so this works whether the deployed env is set to `database` or
; `redis` without needing to match this file to whichever it is. Without a
; worker running at all (there was none), dispatched jobs (e.g.
; GenerateWorkOrdersExcelJob) sit queued forever and their "ready to
; download" notification never fires. --max-time recycles the process
; hourly to shed memory growth; autorestart brings it straight back up.
[program:queue-worker]
command=php artisan queue:work --queue=exports,webhooks,automations,default --sleep=3 --tries=3 --max-time=3600
directory=/var/www/html
autostart=true
autorestart=true
stopwaitsecs=3600
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=30
